# docker/node/supervisord.conf
[supervisord]
nodaemon=true

[program:monitor_ethereum_transaction]
directory=/var/ees
command=yarn monitor:ethereum_transactions
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker_create_in_internal_blockchain]
directory=/var/ees
command=yarn worker:create_in_internal_blockchain
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:monitor_deposit_internal_contract_created]
directory=/var/ees
command=yarn monitor:deposit_internal_contract_created
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:monitor_deposit_internal_contract_redeemed]
directory=/var/ees
command=yarn monitor:deposit_internal_contract_redeemed
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker_redeem_in_external_blockchain]
directory=/var/ees
command=yarn worker:redeem_in_external_blockchain
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker_external_contract_redeem]
directory=/var/ees
command=yarn worker:external_contract_redeem
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:deposit_internal_contract_refunded]
directory=/var/ees
command=yarn monitor:deposit_internal_contract_refunded
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:deposit_internal_contract_burned]
directory=/var/ees
command=yarn monitor:deposit_internal_contract_burned
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:monitor_withdraw_internal_contract_created]
directory=/var/ees
command=yarn monitor:withdraw_internal_contract_created
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:found_withdraw_internal_contract_creation]
directory=/var/ees
command=yarn found:withdraw_internal_contract_creation
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker_create_withdrawal_external_contract]
directory=/var/ees
command=yarn worker:create_withdrawal_external_contract
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker_withdraw_external_contract_redeemed]
directory=/var/ees
command=yarn worker:withdraw_external_contract_redeemed
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:execute_withdraw_internal_contract_redeem]
directory=/var/ees
command=yarn execute:withdraw_internal_contract_redeem
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:monitor_withdraw_internal_contract_redeem]
directory=/var/ees
command=yarn monitor:withdraw_internal_contract_redeem
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:monitor_withdraw_internal_contract_redeem_processed]
directory=/var/ees
command=yarn monitor:withdraw_internal_contract_redeem_processed
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:monitor_external_withdraw_contract_timelock]
directory=/var/ees
command=yarn monitor:external_withdraw_contract_timelock
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:cleanup_temporary_files]
directory=/var/ees
command=rm -r /tmp/yarn*
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=86400
startsecs=0

[program:update_sanctioned_addresses]
directory=/var/ees
command=yarn update:sanctioned-addresses
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=1209600
startsecs=0
